# 开机重启modem
# 1. 用vi创建并编辑文件
vi /etc/init.d/auto-modem

# 2. 按 i 进入编辑模式，粘贴以下内容：
#!/bin/sh /etc/rc.common

START=99
STOP=10

start() {
    sleep 30
    ifup modem
}

stop() {
    ifdown modem
}

# 3. 按 ESC，然后输入 :wq 保存退出

# 4. 设置权限并启用
chmod +x /etc/init.d/auto-modem
/etc/init.d/auto-modem enable
